# ansible/roles/common/templates/healthchecks.timer.j2
[Unit]
Description="Pings healthchecks.io every minute"

[Timer]
OnBootSec=1min
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target